# files/etc/uci-defaults/disable-firewall.sh.orig
#!/bin/sh
echo "Disabling firewall..."
<<<<<<< HEAD
/etc/init.d/firewall disable
=======
chmod 000 /etc/init.d/firewall 
>>>>>>> confine/master

